# services/vhf-ffmpeg.service
[Unit]
Description=VHF Listening Station - FFmpeg Streaming
After=network.target vhf-gnuradio.service
Wants=network.target
Requires=vhf-gnuradio.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/vhf-listening
EnvironmentFile=/home/pi/vhf-listening/.env
ExecStart=/usr/bin/ffmpeg -f f32le -ar 48000 -ac 1 -i /tmp/vhf_audio_pipe -acodec libmp3lame -ab 128k -content_type audio/mpeg -f mp3 icecast://source:${ICECAST_PASSWORD}@${ICECAST_HOST}:${ICECAST_PORT}${ICECAST_MOUNT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target